# debian/cfengine-community.conffiles
/etc/profile.d/cfengine3.sh
/etc/init.d/cfengine3
/etc/default/cfengine3
